# roles/vault/templates/vault.service.j2
[Unit]
Description=HashiCorp Vault Service

[Service]
ExecStart=/opt/SWEAGLE/bin/vaul/vault server -config={{ vault_conf_dir }}/vaultconfig.hcl

[Install]
WantedBy=multi-user.target
