# roles/vault/templates/vaultUnseal.service.j2
[Unit]
Description=HashiCorp Vault Service

[Service]
ExecStart={{ vault_install_dir }}/vaultUnsealStartup.sh

[Install]
WantedBy=multi-user.target
